# scripts/sarracenia.service
[Unit]
Description=Sarracenia Anti-Scraper Tarpit
After=network.target

[Service]
Type=simple
User=sarracenia
Group=sarracenia
WorkingDirectory=/opt/sarracenia

# The binary must be located at /opt/sarracenia/sarracenia
ExecStart=/opt/sarracenia/sarracenia

LimitNOFILE=65536

[Install]
WantedBy=multi-user.target